# keylistend.conf
description "KeyListen daemon"
start on started unity8

exec /home/phablet/keylistend /dev/input/event3 87 dialer-app 88 dialer-app
